# bin/com.macMini.server_update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- THIS FILE HAS TO BE COPIED TO /Library/LaunchDaemons -->
<!-- Then run "launchctl load -w /Library/LaunchDaemons/task.plist" -->
<dict>
    <key>Label</key>
    <string>com.macMini.server_update</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/daniel/dotfiles/bin/server_update.zsh</string>
    </array>
    <key>StartInterval</key>
    <!-- Execute once per hour -->
    <integer>3600</integer>
</dict>
</plist>
